# rootfs/usr/lib/systemd/system/nix-setup.service
[Unit]
Description=Nix store setup
Before=nix-daemon.socket

[Service]
ExecStart=/usr/share/nix/nix-setup
Type=oneshot
