# rootfs/etc/profile.d/nix-env.sh
#!/bin/sh
. /usr/lib/environment.d/nix-daemon.conf
